# emoncms/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Emoncms
# Runs the PHP-FPM daemon
# ==============================================================================

# shellcheck source=../../env-includes/mysql.sh
source /etc/env-includes/mysql.sh
# shellcheck source=../../env-includes/mqtt.sh
source /etc/env-includes/mqtt.sh
# shellcheck source=../../env-includes/etc.sh
source /etc/env-includes/etc.sh

bashio::log.info "Starting PHP-FPM server..."
exec php-fpm7 --nodaemonize
